# systemd/wundio-rfid.service
[Unit]
Description=Wundio RFID Reader Daemon
After=network.target

[Service]
Type=simple
User=wundio
Group=wundio
WorkingDirectory=/opt/wundio/core
EnvironmentFile=/etc/wundio/wundio.env
ExecStart=/opt/wundio/venv/bin/python -c "
import asyncio
from services.rfid import get_rfid_service
async def main():
    s = get_rfid_service()
    s.setup()
    await s.run()
asyncio.run(main())
"
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wundio-rfid

[Install]
WantedBy=multi-user.target
