# systemd/wundio-rfid.service
[Unit]
Description=Wundio RFID Reader Daemon
After=network.target

[Service]
Type=simple
User=wundio
Group=wundio
WorkingDirectory=/opt/wundio/core
EnvironmentFile=/etc/wundio/wundio.env
ExecStart=/opt/wundio/venv/bin/python /opt/wundio/scripts/run-rfid.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wundio-rfid

[Install]
WantedBy=multi-user.target
